# postgres.conf
[program:postgres]
user=root
command=/start-postgis.sh
autorestart=true
stopsignal=INT
