# sshd.conf
[program:sshd]
user=root
command=/usr/sbin/sshd -D
autorestart=true
stopsignal=INT
